# infra/systemd/face-api.service
[Unit]
Description=Face API Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/face-recognition/apps/api
EnvironmentFile=/etc/face/face.env
ExecStart=%h/face-recognition/apps/api/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port ${PORT}
Restart=on-failure

[Install]
WantedBy=multi-user.target

